# deploy/vibe-paradox.service
[Unit]
Description=Vibe Paradox Game Server
After=network.target

[Service]
Type=simple
User=vibe
Group=vibe
WorkingDirectory=/opt/vibe-paradox
ExecStart=/home/vibe/.nvm/versions/node/v20.20.0/bin/node dist/src/server/index.js --port 8080 --db /var/lib/vibe-paradox/world.db
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=ENCRYPTION_KEY=738d5cbdb23fa5ae0cf0ea8a978fb1959b0d45aaff93b8a4823ecf3a7be7f949
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vibe-paradox

[Install]
WantedBy=multi-user.target
